# templates/systemd/portal-gunicorn.service.tmpl
[Unit]
Description=Portal Poczty — aplikacja /admin (Gunicorn)
After=network.target postgresql-${PGDG_MAJOR_VERSION}.service
Requires=postgresql-${PGDG_MAJOR_VERSION}.service

[Service]
Type=simple
User=portal-app
Group=portal-app
WorkingDirectory=/opt/portal-app
RuntimeDirectory=portal-app
RuntimeDirectoryMode=0750

# systemd NIE konsultuje bazy semanage fcontext przy tworzeniu
# RuntimeDirectory= (mimo trwałej reguły dodanej w scripts/vm1/50-portal-app.sh)
# — katalog dostaje zwykłe var_run_t z /run, więc nginx (domena httpd_t) nie
# może się połączyć z gunicorn.sock (AVC denied, "Permission denied" mimo
# 0750/rwx). "+" uruchamia tę jedną linię jako root, mimo User=portal-app
# powyżej — bez tego restorecon sam nie miałby uprawnień.
ExecStartPre=+/usr/sbin/restorecon -R /run/portal-app

# Klucze nigdy nie lądują jako plaintext na dysku — systemd je odszyfrowuje
# (systemd-creds, ewentualnie związane z TPM2) do prywatnej, tylko-w-pamięci
# przestrzeni tego unitu, dostępnej pod $CREDENTIALS_DIRECTORY. Patrz
# scripts/vm1/50-portal-app.sh (generowanie) i portal_app/config.py (odczyt).
LoadCredentialEncrypted=portal-secret-key:/etc/portal/creds/portal-secret-key.cred
LoadCredentialEncrypted=portal-credential-key:/etc/portal/creds/portal-credential-key.cred

ExecStart=/opt/portal-app/venv/bin/gunicorn portal_app.app:app \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind unix:/run/portal-app/gunicorn.sock \
    --access-logfile - --error-logfile -
Restart=on-failure
RestartSec=5

# NoNewPrivileges celowo NIE ustawione: aplikacja wywołuje sudo dla DWÓCH
# ściśle określonych skryptów (apply-branding.sh, apply-tls.sh), a
# NoNewPrivileges=true blokuje jakąkolwiek eskalację przez setuid (sudo pada
# z "no new privileges flag is set"). Realną granicą bezpieczeństwa jest tu
# wąski /etc/sudoers.d/portal-app (tylko te dwa skrypty, root-owned 0700),
# nie ta flaga — identyczny kompromis jak w vm2-provisioning-api.service.
#
# ProtectSystem=full (nie strict): /usr, /boot, /etc pozostają read-only, ale
# /var jest zapisywalny — bez tego sudo'owany helper (dziedziczy namespace
# usługi) nie mógłby zapisać stron błędów do /var/www/errors ani wykonać
# systemctl reload nginx. /etc/portal/tls jawnie odblokowane dla apply-tls.sh.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
# /etc/portal potrzebne, bo apply-branding.sh (sudo, dziedziczy namespace tej
# usługi) zapisuje /etc/portal/roundcube-branding.php, a apply-tls.sh operuje
# na /etc/portal/tls. Bez tego zapis brandingu Roundcube z panelu po cichu
# zawodził (ProtectSystem=full czyni /etc read-only). Katalog i tak jest
# root-owned; realną granicą pozostaje wąski sudoers (tylko te 2 skrypty).
ReadWritePaths=/var/log/portal /var/lib/portal-app /opt/portal-app/portal_app/static /var/www/errors /etc/portal

[Install]
WantedBy=multi-user.target
