# templates/systemd/portal-audit-verify.service.tmpl
[Unit]
Description=Portal Poczty — nocna weryfikacja integralności audit_log (hash-chain)
After=network.target postgresql-${PGDG_MAJOR_VERSION}.service

[Service]
Type=oneshot
User=portal-app
Group=portal-app
WorkingDirectory=/opt/portal-app

LoadCredentialEncrypted=portal-secret-key:/etc/portal/creds/portal-secret-key.cred
LoadCredentialEncrypted=portal-credential-key:/etc/portal/creds/portal-credential-key.cred

ExecStart=/opt/portal-app/venv/bin/python -m portal_app.workers.audit_verify_main

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
