# templates/systemd/portal-worker.service.tmpl
[Unit]
Description=Portal Poczty — worker kolejki zadań (sync/provision)
After=network.target postgresql-${PGDG_MAJOR_VERSION}.service portal-gunicorn.service
Requires=postgresql-${PGDG_MAJOR_VERSION}.service

[Service]
Type=simple
User=portal-app
Group=portal-app
WorkingDirectory=/opt/portal-app
# Oba katalogi tworzone przez systemd PRZED montowaniem namespace (ProtectSystem):
#  portal-app    — współdzielony z gunicornem (gniazdo gunicorn.sock)
#  portal-import — tmpfs na strumieniowe rozpakowanie archiwów importu (config.py
#                  IMPORT_TMP_DIR). KLUCZOWE: po reboocie/`dnf update` /run (tmpfs)
#                  jest czyszczony; gdyby /run/portal-import było tylko w
#                  ReadWritePaths (a nie RuntimeDirectory), worker padłby na
#                  226/NAMESPACE ("Failed to set up mount namespacing: ... No such
#                  file or directory"). RuntimeDirectory gwarantuje ich istnienie
#                  przy każdym starcie — także pierwszym po reboocie.
RuntimeDirectory=portal-app portal-import
RuntimeDirectoryMode=0750
# Preserve=yes: katalogi przeżywają restart którejkolwiek z usług (restart workera
# nie skasuje wtedy gniazda gunicorna → koniec z 502 po restarcie workera).
RuntimeDirectoryPreserve=yes

LoadCredentialEncrypted=portal-secret-key:/etc/portal/creds/portal-secret-key.cred
LoadCredentialEncrypted=portal-credential-key:/etc/portal/creds/portal-credential-key.cred

ExecStart=/opt/portal-app/venv/bin/python -m portal_app.workers.worker_main
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# /run/portal-app i /run/portal-import są RuntimeDirectory (automatycznie
# zapisywalne, wyłączone z ProtectSystem) — nie muszą być w ReadWritePaths.
ReadWritePaths=/var/log/portal /var/lib/portal-app

[Install]
WantedBy=multi-user.target
