# templates/systemd/portal-worker.service.tmpl
[Unit]
Description=Portal Poczty — worker kolejki zadań (sync/provision)
After=network.target postgresql-${PGDG_MAJOR_VERSION}.service portal-gunicorn.service
Requires=postgresql-${PGDG_MAJOR_VERSION}.service

[Service]
Type=simple
User=portal-app
Group=portal-app
WorkingDirectory=/opt/portal-app
# Oba katalogi tworzone przez systemd PRZED montowaniem namespace (ProtectSystem):
#  portal-app    — współdzielony z gunicornem (gniazdo gunicorn.sock)
#  portal-import — tmpfs na strumieniowe rozpakowanie archiwów importu (config.py
#                  IMPORT_TMP_DIR). KLUCZOWE: po reboocie/`dnf update` /run (tmpfs)
#                  jest czyszczony; gdyby /run/portal-import było tylko w
#                  ReadWritePaths (a nie RuntimeDirectory), worker padłby na
#                  226/NAMESPACE ("Failed to set up mount namespacing: ... No such
#                  file or directory"). RuntimeDirectory gwarantuje ich istnienie
#                  przy każdym starcie — także pierwszym po reboocie.
RuntimeDirectory=portal-app portal-import
RuntimeDirectoryMode=0750
# Preserve=yes: katalogi przeżywają restart którejkolwiek z usług (restart workera
# nie skasuje wtedy gniazda gunicorna → koniec z 502 po restarcie workera).
RuntimeDirectoryPreserve=yes

LoadCredentialEncrypted=portal-secret-key:/etc/portal/creds/portal-secret-key.cred
LoadCredentialEncrypted=portal-credential-key:/etc/portal/creds/portal-credential-key.cred

ExecStart=/opt/portal-app/venv/bin/python -m portal_app.workers.worker_main
Restart=on-failure
RestartSec=5

# NoNewPrivileges celowo NIE ustawione: worker wywołuje sudo dla wąsko
# zdefiniowanych helperów aktualizacji systemu (apply-system-update.sh
# backup|update|health|reboot-check) — aktualizacje idą w tle przez workera, bo
# długi dnf zabiłby żądanie web (gunicorn --timeout). NoNewPrivileges=true
# blokuje jakąkolwiek eskalację przez setuid (sudo pada z "no new privileges").
# Realną granicą jest wąski /etc/sudoers.d/portal-app (tylko te helpery,
# root-owned 0700) — identyczny kompromis jak w portal-gunicorn.service.
# ProtectSystem=full (nie strict): /usr,/boot,/etc read-only, ale /var
# zapisywalny — helper backupu pisze /var/lib/portal-config-backups, a dnf i tak
# escape'uje sandbox przez systemd-run (transient unit na hoście).
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
# /run/portal-app i /run/portal-import są RuntimeDirectory (automatycznie
# zapisywalne). /var jest zapisywalny pod ProtectSystem=full — kopie konfiguracji
# trafiają do /var/lib/portal-config-backups.
ReadWritePaths=/var/log/portal /var/lib/portal-app /var/lib/portal-config-backups

[Install]
WantedBy=multi-user.target
