# templates/systemd/vm2-provisioning-api.service.tmpl
[Unit]
Description=VM2 Provisioning API (mTLS, wywoływane wyłącznie z VM1)
After=network-online.target postgresql-${PGDG_MAJOR_VERSION}.service
Wants=network-online.target

[Service]
Type=simple
User=vm2-api
Group=vm2-api
WorkingDirectory=/opt/vm2-api
Environment=VM2_API_ALLOWED_CLIENT_IP=${VM1_IP}
Environment=VM2_API_PORT=${VM2_API_PORT}
ExecStart=/opt/vm2-api/venv/bin/python -m vm2_api.main
Restart=on-failure
RestartSec=5

# NoNewPrivileges NIE jest ustawione celowo: usługa musi wywoływać `sudo`
# dla ściśle ograniczonego zestawu komend uprzywilejowanych (dnf update,
# systemctl reboot, freshclam, clamdscan) — sudo wymaga setuid, co
# NoNewPrivileges by zablokowało. Rzeczywistą granicą bezpieczeństwa jest tu
# wąski /etc/sudoers.d/vm2-api, nie flaga systemd. Patrz
# docs/technical/architecture.md.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# /var/mail/vhosts MUSI być zapisywalny w namespace usługi: `doveadm quota get`
# (backend quota-maildir) tworzy/aktualizuje plik `maildirsize` w Maildirze, by
# policzyć zajętość — pod ProtectSystem=strict katalog jest read-only, doveadm
# zwraca wtedy "error" i panel pokazuje 0 B zajętości skrzynki. Zapis dotyczy
# tylko pliku licznika quoty; samych maili nie modyfikujemy.
ReadWritePaths=/opt/vm2-api /var/mail/vhosts

[Install]
WantedBy=multi-user.target
